# deploy/systemd/prometheus.service
[Unit]
Description=Prometheus Server
After=network.target

[Service]
User=prometheus
Group=prometheus
WorkingDirectory=/opt/spotify-optimizer/deploy/monitoring
ExecStart=/usr/bin/prometheus --config.file=/opt/spotify-optimizer/deploy/monitoring/prometheus.yml --storage.tsdb.path=/var/lib/prometheus
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/prometheus

[Install]
WantedBy=multi-user.target
