# deploy/systemd/spotify-backend.service
[Unit]
Description=Spotify Optimizer Backend
After=network.target redis-server.service

[Service]
User=www-data
WorkingDirectory=/opt/spotify-optimizer
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/opt/spotify-optimizer/backend/.env
ExecStart=/opt/spotify-optimizer/backend/.venv/bin/uvicorn backend.app:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
UMask=027
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/spotify-optimizer/backend/data /opt/spotify-optimizer/backend/cache /opt/spotify-optimizer/backend/logs /opt/spotify-optimizer/backend/models
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
